# apps/camera-agent/deploy/systemd/proof-agent.service
[Unit]
Description=Proof of Reality capture agent (B2B)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/proof-agent
EnvironmentFile=/opt/proof-agent/.env
ExecStart=/usr/local/bin/proof-agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
